# deploy/nfe-web.service
[Unit]
Description=Automação NFe Web (FastAPI)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=nfe-web
Group=nfe-web
WorkingDirectory=/opt/nfe-web
Environment=PYTHONPATH=/opt/nfe-web
EnvironmentFile=/opt/nfe-web/.env
ExecStart=/opt/nfe-web/venv/bin/python -m uvicorn backend.app.main:app --host 0.0.0.0 --port 8090
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
